# nginx-lua-edge-service/data/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:amplify-agent]
command=/bin/bash -c "service amplify-agent start"

[program:nginx]
command=/bin/bash -c "/usr/sbin/nginx -c \"/data/etc/nginx/nginx.conf\" -g \"daemon off;\""